# conf/ns_supervisor.conf
[program:ns]
command=/opt/myenv/smartlamp/bin/python /usr/share/smartlamp_core/start_server.py
startsecs = 0
autorestart = true
stdout_logfile=/var/log/ns/ns.log