# conf/ns_supervisor.conf
[program:ns]
command=/opt/myenv/smartlamp/bin/python /usr/share/smartlamp_core/start_server.py
stdout_logfile=/var/log/ns/ns_supervisor.log